# salt/etc/udev/rules.d/99-discharge.rules
SUBSYSTEM=="power_supply", ATTR{status}=="Discharging", RUN+="{{ grains['cbi_home'] }}/bin/discharging %E{POWER_SUPPLY_CHARGE_NOW}"