# salt/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin {{ pillar['desktop_user'] }} --noclear %I 38400 linux
Type=simple